# service/src/main/daemon/red5.service
[Unit]
Description=Red5

[Service]
Type=simple
ExecStart=/etc/init.d/red5

[Install]
WantedBy=multi-user.target
